# infras/config/supervisord/app.conf
[program:app]
command=/usr/local/bin/app
numprocs=1
autostart=true
autorestart=false